# dot_config/systemd/user/clipman.service
[Service]
ExecStart=/usr/bin/wl-paste -t text --watch clipman store

[Install]
WantedBy=sway-session.target hyprland-session.target
