# dot_config/systemd/user/kanshi.service
[Service]
ExecStart=/usr/bin/kanshi
Restart=always

[Install]
WantedBy=sway-session.target hyprland-session.target
