# dot_config/systemd/user/swayosd.service
[Service]
ExecStart=/usr/bin/swayosd
Restart=always

[Install]
WantedBy=sway-session.target hyprland-session.target
