# dot_config/systemd/user/waybar.service
[Service]
ExecStart=/usr/bin/waybar -b bar-0
Restart=always

[Install]
WantedBy=sway-session.target hyprland-session.target
